# tools/scheduling/guide.m3c.skillctl-sweep.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- SPEC-0247/SPEC-0251 §5: daily skillctl revocation sweep (macOS launchd).
     EDIT the skillctl path below for your machine, then:
       cp this ~/Library/LaunchAgents/ && launchctl load -w ~/Library/LaunchAgents/guide.m3c.skillctl-sweep.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>guide.m3c.skillctl-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/skillctl</string>
    <string>verify</string>
    <string>--all</string>
    <string>--quarantine</string>
    <string>--json</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/skillctl-sweep.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/skillctl-sweep.err.log</string>
</dict>
</plist>
